# DataAInsint/backend/dataainsight.service
[Unit]
Description=DataAInsight Data Exploration Tool
After=network.target

[Service]
# 服务类型
Type=forking

# 用户和组（根据实际情况修改）
User=www-data
Group=www-data

# 工作目录
WorkingDirectory=/path/to/dataAInsight/backend

# 启动命令
ExecStart=/path/to/dataAInsight/backend/start_prod_linux.sh start
ExecStop=/path/to/dataAInsight/backend/start_prod_linux.sh stop
ExecReload=/path/to/dataAInsight/backend/start_prod_linux.sh restart

# PID文件
PIDFile=/tmp/dataainsight.pid

# 内存限制配置
MemoryLimit=2G
MemoryAccounting=yes
MemoryMax=2G

# CPU限制（可选）
CPUQuota=200%

# 文件描述符限制
LimitNOFILE=65536

# 进程数限制
LimitNPROC=4096

# 虚拟内存限制
LimitAS=2147483648

# 重启策略
Restart=always
RestartSec=10

# 环境变量
Environment=PYTHONPATH=/path/to/dataAInsight/backend
Environment=PYTHONUNBUFFERED=1

# 日志配置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dataainsight

[Install]
WantedBy=multi-user.target